# systemd/openclaw-gateway.service.d/10-auto-mihomo.conf
# Auto-Mihomo drop-in overlay for openclaw-gateway.
# The base unit is owned and updated by OpenClaw (openclaw onboard --install-daemon).
# Only this file is managed by auto-mihomo.

[Unit]
After=network-online.target mihomo.service
Wants=network-online.target mihomo.service

[Service]
# Proxy environment (HTTP_PROXY / HTTPS_PROXY / NO_PROXY) — refreshed by update_sub.sh
EnvironmentFile=/etc/auto-mihomo/proxy.env
# Ensure nvm-managed Node.js is on PATH (systemd starts with minimal PATH)
Environment=PATH=__NODE_DIR__:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
# Preload proxy-bootstrap so undici uses EnvHttpProxyAgent before the first fetch
# (openclaw/openclaw#42320 has not merged yet)
Environment=NODE_OPTIONS=-r __PROJECT_DIR__/scripts/proxy-bootstrap.cjs
# Update subscription before gateway starts; the - prefix tolerates failure
ExecStartPre=-bash __PROJECT_DIR__/scripts/update_sub.sh
