# systemd/openclaw-gateway.service
[Unit]
Description=OpenClaw Gateway
After=network-online.target mihomo.service
Wants=network-online.target mihomo.service

[Service]
Type=simple
User=__USER__
WorkingDirectory=__HOME__/.openclaw

# 从 auto-mihomo 统一管理的代理配置读取 (端口变更自动生效)
EnvironmentFile=/etc/auto-mihomo/proxy.env

# Node.js 内置 fetch (undici) 不读取 HTTP_PROXY 环境变量,
# 通过 proxy-bootstrap.cjs 预加载 undici ProxyAgent (零外部依赖)
Environment="HOME=__HOME__"
Environment="NODE_OPTIONS=-r __PROJECT_DIR__/scripts/proxy-bootstrap.cjs"

ExecStart=__PROJECT_DIR__/scripts/start_openclaw_with_proxy.sh
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
