# systemd/openclaw-gateway.service
[Unit]
Description=OpenClaw Gateway
After=network-online.target mihomo.service
Wants=network-online.target mihomo.service

[Service]
# 本模板面向 system 级常驻部署；install.sh / upgrade.sh 会额外 mask
# ${HOME}/.config/systemd/user/openclaw-gateway.service，避免官方 CLI fallback 到
# systemctl --user 时拉起第二个 gateway 并抢占 18789 端口。
Type=simple
User=__USER__
WorkingDirectory=__HOME__/.openclaw

# 从 auto-mihomo 统一管理的代理配置读取 (端口变更自动生效)
EnvironmentFile=/etc/auto-mihomo/proxy.env

Environment="HOME=__HOME__"
Environment="TMPDIR=/tmp"
Environment="PATH=__NODE_DIR__:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="OPENCLAW_GATEWAY_PORT=18789"
Environment="OPENCLAW_SYSTEMD_UNIT=openclaw-gateway.service"
Environment="OPENCLAW_WINDOWS_TASK_NAME=OpenClaw Gateway"
Environment="OPENCLAW_SERVICE_MARKER=openclaw"
Environment="OPENCLAW_SERVICE_KIND=gateway"
Environment="OPENCLAW_SERVICE_VERSION=__OPENCLAW_SERVICE_VERSION__"
Environment="OPENCLAW_STATE_DIR=__HOME__/.openclaw"
Environment="OPENCLAW_CONFIG_PATH=__HOME__/.openclaw/openclaw.json"

ExecStart=__PROJECT_DIR__/scripts/start_openclaw_with_proxy.sh
Restart=always
RestartSec=5
TimeoutStopSec=30
TimeoutStartSec=30
SuccessExitStatus=0 143
KillMode=control-group

[Install]
WantedBy=multi-user.target
